# image-builder/postinstall/kek/overlay/etc/supervisor/conf.d/gpd4303s-ioc.conf
[program:gpd4303s-ioc]
command=/opt/GPD4303S-ioc/python/GPD4303S-ioc.py
stdout_logfile=/var/log/gpd4303s-ioc.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
stderr_logfile=/var/log/gpd4303s-ioc.log
stderr_ogfile_maxbytes=1MB
stderr_logfile_backups=5
